# apps/golang/rootfs/etc/profile.d/10-go-env.sh
# NiceOS Golang environment defaults
# SPDX-License-Identifier: Apache-2.0

export GOPATH="${GOPATH:-/go}"
export GOCACHE="${GOCACHE:-/go/.cache}"
export GOMODCACHE="${GOMODCACHE:-/go/pkg/mod}"
export GOROOT="${GOROOT:-/opt/bitnami/go}"
export HOME="${HOME:-/go}"
